# systemd/dvbondemand-dev.path
[Unit]
Description=Gets audio/video services from several dvb cards.
#Documentation=man:dvbondemand(1)

[Path]
PathChanged=/dev/dvb/
PathChanged=/etc/default/dvbondemand
Unit=dvbondemand-dev.service

[Install]
WantedBy=multi-user.target
